# .config/autostart/pulseaudio.desktop
[Desktop Entry]
Name=Pulseaudio
Comment=Pulseaudio
Exec=pulseaudio
Terminal=false
Type=Application
